# packaging/usr/lib/systemd/user/crush-claw.service
[Unit]
Description=Crush Claw — esp-claw Desktop AI Agent
Documentation=https://github.com/espressif/esp-claw
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/esp-claw-desktop --foreground --pid-file %h/.crush-claw/agent.pid
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=crush-claw
Environment=CRUSH_CLAW_DATA_DIR=%h/.crush-claw

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/.crush-claw
ReadOnlyPaths=/usr/share/crush-claw

[Install]
WantedBy=default.target
